# Mone/deploy/mone-songbook.service
[Unit]
Description=MONE Interactive Songbook Viewer
After=network.target

[Service]
Type=simple
User=moneapp
Group=moneapp
WorkingDirectory=/opt/mone-songbook
EnvironmentFile=/etc/mone-songbook.env
ExecStart=/usr/bin/node /opt/mone-songbook/server.js
Restart=always
RestartSec=2
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
